# config/systemd/vncserver.service
# NOTE: For now this is only a proof of concept and should not be used

[Unit]
Description=VNC server
After=network.target

[Service]
Type=forking
WorkingDirectory=/home/simon/shared/code/neopult
ExecStart=/home/simon/shared/code/neopult/scripts/vnc-start.sh -d

[Install]
WantedBy=multi-user.target
